# examples/systemd/pyanchor.service
# /etc/systemd/system/pyanchor.service
#
# pyanchor sidecar systemd unit — production-hardened template.
#
# This is the same hardening block from docs/PRODUCTION-HARDENING.md
# pulled into a standalone file so you can `cp` it directly without
# copy/pasting from a markdown doc. Pair with pyanchor.env.example
# (in this same directory) for the matching EnvironmentFile.
#
# Quick start (after editing both files for your paths):
#   sudo cp pyanchor.service /etc/systemd/system/
#   sudo cp pyanchor.env.example /etc/pyanchor.env
#   sudo nano /etc/pyanchor.env       # fill in tokens / paths
#   sudo chmod 600 /etc/pyanchor.env  # secrets, root-only
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pyanchor
#   sudo systemctl status pyanchor
#   curl http://127.0.0.1:3010/healthz   # liveness
#   curl http://127.0.0.1:3010/readyz    # readiness (200 only when fully configured)
#
# Verify the sandbox is active (lower is better):
#   systemd-analyze security pyanchor

[Unit]
Description=pyanchor sidecar (AI live-edit for Next.js / Vite)
Documentation=https://github.com/pyanchor/pyanchor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pyanchor
Group=pyanchor

# Working directory + env. Adjust paths for your install layout.
WorkingDirectory=/opt/pyanchor
EnvironmentFile=/etc/pyanchor.env
ExecStart=/usr/bin/node /opt/pyanchor/dist/server.cjs

# Restart policy. RestartSec=5s prevents tight-loop crash restarts
# from drowning the journal when the agent CLI is missing.
Restart=on-failure
RestartSec=5s

# ─── Sandbox (free, fail-closed, survives code bugs) ───
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Optional further hardening — DISABLED BY DEFAULT because it breaks
# Node/V8's JIT on some builds and may break agent CLIs that use
# native JIT runtimes. Enable only after verifying pyanchor + your
# selected agent + your output mode all work under systemd:
#   MemoryDenyWriteExecute=true
# See systemd.exec(5) and https://github.com/nodejs/node/issues for
# current V8 W^X status on your Node version.

# Allow writes ONLY to the dirs the worker needs.
#   /var/lib/pyanchor/workspace — scratch space the agent mutates
#   /var/lib/pyanchor/state     — state.json + audit.jsonl
#   /srv/myapp                  — PYANCHOR_APP_DIR (apply mode only;
#                                 remove this line for pr mode, which
#                                 writes via git push instead of rsync)
# If PYANCHOR_APP_DIR differs from /srv/myapp, update the third path
# or the sidecar will get EROFS during rsync.
ReadWritePaths=/var/lib/pyanchor/workspace /var/lib/pyanchor/state /srv/myapp

# Network: pyanchor binds 127.0.0.1 by default (controlled by
# PYANCHOR_HOST) and your reverse proxy handles the public face.
# Do NOT block outbound traffic here — agent CLIs call LLM providers,
# PR mode calls the GitHub API, webhooks call Slack/Discord. If you
# need site-specific egress restrictions, add an IPAddressDeny /
# IPAddressAllow policy tailored to your provider's egress ranges
# rather than a blanket loopback-only deny.

# Resource limits. The worker can spawn agents that allocate large
# context windows; bump MemoryMax if you see OOM kills under load.
LimitNOFILE=4096
MemoryMax=2G
TasksMax=256

[Install]
WantedBy=multi-user.target
